# build/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=/usr/local/bin/node /app/server/index.mjs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:caddy]
command=/usr/bin/caddy run --adapter caddyfile --config /etc/caddy/Caddyfile
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
